# board/common/rootfs_overlay/etc/systemd/system/carui.service
[Unit]
Description=carui

[Service]
Type=simple
ExecStart=/usr/bin/CarUI

[Install]
WantedBy=basic.target
