# board/common/rootfs_overlay/etc/systemd/system/pulseaudio.service
[Unit]
Description=pulseaudio

[Service]
Type=notify
ExecStart=/usr/bin/pulseaudio --daemonize=no

[Install]
WantedBy=basic.target
